# deploy/deleqate.service
[Unit]
Description=Deleqate Django backend (gunicorn)
After=network.target

[Service]
User=opc
Group=opc
WorkingDirectory=/home/opc/deleqate/backend
Environment="PATH=/home/opc/deleqate/backend/venv/bin"
# gunicorn serves the Django WSGI app on a local port; nginx proxies to it.
ExecStart=/home/opc/deleqate/backend/venv/bin/gunicorn config.wsgi:application \
    --bind 127.0.0.1:5051 \
    --workers 3 \
    --timeout 120
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
